# src/selfheal/systemd/selfheal.target
[Unit]
Description=SelfHeal Target
StopWhenUnneeded=yes
